# DEBIAN/postinst
#!/bin/sh
systemctl daemon-reload
systemctl enable watchlog-agent
systemctl start watchlog-agent